# hipache.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="hipache">

  <service name="site/hipache" type="service" version="3">

    <create_default_instance enabled="true"/>

    <single_instance/>

    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default"/>
    </dependency>

    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local"/>
    </dependency>

    <method_context>
      <method_credential user="hipache" group="hipache"/>
    </method_context>

    <exec_method type="method" name="start" exec="/opt/hipache/bin/hipache -c %{config_file}" timeout_seconds="60">
      <method_context>
        <method_credential user='hipache' group='hipache' privileges='basic,net_privaddr'/>
        <method_environment>
          <envvar name="PATH" value="/opt/local/bin/" />
        </method_environment>
      </method_context>
    </exec_method>

    <exec_method type="method" name="stop" exec=":kill" timeout_seconds="60"/>

    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="child"/>


      <propval name="ignore_error" type="astring" value="core,signal"/>
    </property_group>

    <property_group name="application" type="application">
      <propval name="config_file" type="astring" value="/etc/hipache.json"/>
    </property_group>


    <stability value="Evolving"/>

    <template>
      <common_name>
        <loctext xml:lang="C">
          Hipache Proxy
        </loctext>
      </common_name>
    </template>

  </service>

</service_bundle>


