# upstart.conf
# hipache - distributed HTTP and websocket proxy

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 15 5

script
  hipache --config /etc/hipache.json
end script